# scripts/drive/systemd/blink-drive-test.service
[Unit]
Description=Blink Drive Test Service
After=network.target

[Service]
Type=oneshot
ExecStart=/bin/bash -c 'echo "Starting USB gadget..." && /opt/blink-sync-brain/scripts/drive/start_storage_mode.sh && echo "USB gadget started successfully"'
RemainAfterExit=yes
User=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
